# launchd/io.mutagen.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>io.mutagen.daemon</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/sh</string>
        <string>-c</string>
        <string>export PATH="$HOME/.nix-profile/bin:/nix/var/nix/profiles/default/bin:$PATH"; mutagen daemon start</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/tmp/mutagen-daemon.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/mutagen-daemon.err</string>
</dict>
</plist>
